# root/etc/s6-overlay/s6-rc.d/svc-mod-universal-code-server-code-tunnel/run
#!/usr/bin/with-contenv bash
echo "running service"
/usr/bin/code tunnel --accept-server-license-terms