# root/etc/s6-overlay/s6-rc.d/init-mod-universal-code-server-code-tunnel/up
/etc/s6-overlay/s6-rc.d/init-mod-universal-code-server-code-tunnel/run
