# root/etc/services.d/code-server-code-tunnel/run
#!/usr/bin/with-contenv bash

/usr/bin/code tunnel --accept-server-license-terms